# archivedir.service
[Unit]
Description=Archive directories per configuration

[Service]
Type=oneshot
WorkingDirectory=/tmp
ExecStart=/usr/bin/archivedir-runner
